# crontab
# log the weather every minute
* * * * *	/home/rcf/bin/log-weather
