# crontab
# log the weather every minute
* * * * *	/home/weather/bin/weather-log
